# puppet/modules/blueblazer/templates/blueblazer.conf
[program:blueblazer]
command=/usr/bin/blueblazer-http
directory=/var/www/blueblazer/bin
redirect_stderr=true
