# puppet/modules/blueblazer/templates/blueblazer.conf
[program:blueblazer]
command=/var/www/blueblazer/bin/blueblazer-http
directory=/var/www/blueblazer/bin
redirect_stderr=true
